# systemd/onair.service
[Unit]
Description=A description of what your service does.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/on-air